# systemd/for-server.service
[Unit]
Description=For Automation Framework Server
After=network.target
Documentation=https://github.com/diceone/for-IT

[Service]
Type=simple
User=for
Group=for
ExecStartPre=/bin/mkdir -p /var/log/for
ExecStartPre=/bin/chown for:for /var/log/for
ExecStartPre=/bin/chmod 755 /var/log/for
ExecStart=/usr/local/bin/for-server --addr :8080 --playbook-dir /etc/for/environments --log /var/log/for/server.log
Restart=always
RestartSec=10
WorkingDirectory=/etc/for
StandardOutput=append:/var/log/for/server.log
StandardError=append:/var/log/for/server.error.log

# Security settings
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
