# systemd/for-client.service
[Unit]
Description=For Client Service
After=network.target

[Service]
Type=simple
User=for
Group=for
RuntimeDirectory=for
RuntimeDirectoryMode=0755
LogsDirectory=for
LogsDirectoryMode=0755
WorkingDirectory=/etc/for
ExecStartPre=/bin/mkdir -p /var/log/for
ExecStartPre=/bin/chown for:for /var/log/for
ExecStartPre=/bin/chmod 755 /var/log/for
ExecStart=/usr/local/bin/for-client -server localhost:8080 -customer customer1 -environment production
Restart=always
RestartSec=10
StandardOutput=append:/var/log/for/client.log
StandardError=append:/var/log/for/client.error.log

# Security settings
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
